# airootfs/etc/systemd/system/aero-firstboot.service
[Unit]
Description=Aero Linux First Boot Setup
Documentation=https://github.com/aero-linux
After=network-online.target graphical.target
Wants=network-online.target
# Run only on an Aero-installed system where first-boot hasn't completed yet
ConditionPathExists=/etc/aero-installed
ConditionPathExists=!/etc/aero-firstboot-complete

[Service]
Type=oneshot
ExecStart=/usr/share/aero/scripts/first-boot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target